# Raspberry/wifi.service
[Unit]
After=network.target

[Service]
ExecStart=/home/root/connect.sh

[Install]
WantedBy=default.target